# example/manage-ssl-cert.timer
[Unit]
Description=Keep the SSL certificate up to date

[Timer]
OnBootSec=6h
OnUnitActiveSec=6h

[Install]
WantedBy=timers.target
